# contrib/vultrack-agent.service
[Unit]
Description=VulTrack Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vultrack-agent
Group=vultrack-agent
ExecStart=/usr/local/bin/vultrack-agent daemon
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/vultrack-agent /var/lib/vultrack-agent
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
