# ansible/roles/vultrack_agent/templates/vultrack-agent.service.j2
# {{ ansible_managed }}
[Unit]
Description=VulTrack Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vultrack-agent
Group=vultrack-agent
ExecStart={{ vultrack_agent_bin_path }} daemon
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ vultrack_agent_config_dir }} /var/lib/vultrack-agent
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
